# test/container/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:postfix]
command=/usr/lib/postfix/sbin/master -c /etc/postfix -d

[program:dovecot]
command=/usr/sbin/dovecot -F

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
